# systemd/arrivals-led.service
[Unit]
Description=Arrivals LED Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/arrivals-led
ExecStart=%h/arrivals-led/venv/bin/python %h/arrivals-led/arrivals.py "arrivals --json tfl --station 910GSHRDHST --platform 2"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
